# deploy/oquv-web.service
[Unit]
Description=Oquv Markaz - Web panel (uvicorn)
After=network.target

[Service]
WorkingDirectory=/root/oquv_markaz
ExecStart=/root/oquv_markaz/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
