# deploy/oquv-reminders.service
[Unit]
Description=Oquv Markaz - kunlik eslatmalar (to'lov + tug'ilgan kun)

[Service]
Type=oneshot
WorkingDirectory=/root/oquv_markaz
ExecStart=/root/oquv_markaz/.venv/bin/python -m app.reminders
