# deploy/oquv-reminders.timer
[Unit]
Description=Har kuni 09:00 da eslatmalarni yuborish

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
